# deploy/supervisord.conf
[program:vision-mcp]
command=/srv/vision-mcp/.venv/bin/python -m deploy.launcher
directory=/srv/vision-mcp
user=ubuntu
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10
stdout_logfile=/var/log/vision-mcp/stdout.log
stderr_logfile=/var/log/vision-mcp/stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=
    ARK_BASE_URL="__ARK_BASE_URL__",
    ARK_API_KEY="__ARK_API_KEY__",
    ARK_MODEL="__ARK_MODEL__",
    VISION_BEARER_TOKEN="__VISION_BEARER_TOKEN__",
    HOST="127.0.0.1",
    PORT="8100"